# bridge/launchd/com.local.claude-quota-bridge.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd job for the quota bridge (serves the LAN, makes no outbound calls).

  Setup:
    1. cp this file to ~/Library/LaunchAgents/, dropping the .example suffix.
    2. Replace /Users/you/esp32-claude-quota with your actual repo path.
    3. launchctl load ~/Library/LaunchAgents/com.local.claude-quota-bridge.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.local.claude-quota-bridge</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/you/esp32-claude-quota/bridge/claude-quota-bridge</string>
    <string>--no-tokens</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/tmp/claude-quota-bridge.out</string>
  <key>StandardErrorPath</key><string>/tmp/claude-quota-bridge.err</string>
</dict>
</plist>
